# config/zshWeather/weather-fetcher.service
# /etc/systemd/system/weather-fetcher.service
[Unit]
Description=Fetch weather data

[Service]
ExecStart=/bin/zshWeather
